# root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv bash

pushd /opt/piavpn || exit
exec /opt/piavpn/up.sh
