# root/etc/services.d/taillogs/run
#!/usr/bin/with-contenv bash

exec tail -F -n0 /config/deluged.log /opt/piavpn-manual/debug_info 2> /dev/null
